# supervisord.conf
[supervisord]
nodaemon=true

[program:setup-script]
command = bash -c "/setup.sh"
startsecs = 0
autorestart = false
startretries = 1

[program:sshd]
command=/usr/sbin/sshd -D
